# raspberry_pi/systemd/irig-sender.service
[Unit]
Description=IRIG-H Timecode Sender
# Ensure chrony is running before we start polling it
After=multi-user.target chronyd.service
Wants=chronyd.service

[Service]
Type=simple
# Wait until chrony offset < 50 ms before starting transmission.
# Prevents clock steps from introducing discontinuities in the IRIG timecode.
# chronyc waitsync args: max-tries=0 (unlimited), delta=0 (ignore), offset=0.05 s, interval=10 s
ExecStartPre=/usr/bin/chronyc waitsync 0 0 0.05
# Pin configuration is set by install.sh (e.g., ./install.sh -p 17 -n 27).
# The ExecStart line below is the default; install.sh substitutes it.
ExecStart=/usr/local/bin/irig_sender -p 9 -n -1
Restart=always
RestartSec=10
User=root
Nice=-20
StandardOutput=journal
StandardError=journal

# Create a working directory for output files
WorkingDirectory=/var/log/irig-sender
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Ensure the service stops gracefully
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target