# raspberry_pi/systemd/irig-sender.service
[Unit]
Description=IRIG-H Timecode Sender
# Ensure chrony is running before we start polling it
After=multi-user.target chronyd.service
Wants=chronyd.service

[Service]
Type=simple
# Pin configuration is set by install.sh (e.g., ./install.sh -p 17 -n 27).
# The ExecStart line below is the default; install.sh substitutes it.
ExecStart=/usr/local/bin/irig_sender -p 9 -n -1
Restart=always
RestartSec=10
User=root
Nice=-20
StandardOutput=journal
StandardError=journal

# Create a working directory for output files
WorkingDirectory=/var/log/irig-sender
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Ensure the service stops gracefully
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target